# claude-agent.service
[Unit]
Description=Claude AI Agent - Offline AI Assistant
After=network.target ollama.service
Requires=ollama.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/claude-agent
ExecStart=/usr/bin/node /root/claude-agent/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-agent

# Environment
Environment=NODE_ENV=production
Environment=PORT=7778
Environment=OLLAMA_HOST=http://localhost:11434
Environment=OLLAMA_MODEL=llama3
Environment=DB_PATH=./database/agent.db

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target